# template.service
[Unit]
After={after}
{unit-options}
[Service]
ExecStart={command}
Restart=always
{service-options}
[Install]
WantedBy=multi-user.target
